# scripts/pi/chromium-kiosk.service
[Unit]
Description=Chromium kiosk for FamilyBoard
After=graphical-session.target

[Service]
Type=simple
# Hide the cursor before launching the browser — unclutter watches the X root window.
ExecStartPre=/usr/bin/unclutter -idle 0 -root -b
# On-screen keyboard for third-party pages (Google login) our React keyboard
# can't reach. Best-effort (leading '-'): never block the kiosk if onboard is
# absent. The systemd --user session already provides a D-Bus session bus.
ExecStartPre=-/bin/bash -c 'command -v onboard >/dev/null 2>&1 && { gsettings set org.onboard auto-show.enabled true; gsettings set org.onboard.window.docking enabled true; gsettings set org.onboard start-minimized true; setsid onboard >/dev/null 2>&1 & } || true'
ExecStartPre=/bin/bash -c 'until curl -fsS http://localhost:3000/api/health > /dev/null; do sleep 2; done'
ExecStart=/usr/bin/chromium-browser \
  --kiosk \
  --noerrdialogs \
  --disable-infobars \
  --disable-session-crashed-bubble \
  --disable-translate \
  --no-first-run \
  --start-fullscreen \
  --check-for-update-interval=31536000 \
  --autoplay-policy=no-user-gesture-required \
  --overscroll-history-navigation=0 \
  --force-renderer-accessibility \
  http://localhost:3000
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
