# scripts/pi/familyboard-updater.path
[Unit]
Description=Watch for a user-triggered FamilyBoard update request

[Path]
# The wall UI (Settings -> "check for updates") writes this flag file into the
# bind-mounted data dir; the app never needs the Docker socket. The updater
# removes the flag after a successful run.
PathExists=/opt/familyboard/data/update-request
Unit=familyboard-updater.service

[Install]
WantedBy=paths.target
